# systemd/goimapnotify.service.template
[Unit]
Description=IMAP IDLE mail notifier (goimapnotify)
After=network-online.target

[Service]
ExecStart=/usr/bin/goimapnotify -wait 10 -conf ${REPO_DIR}/config/goimapnotify/config.json
Restart=on-failure
RestartSec=15s

[Install]
WantedBy=default.target
